# motor_control/motor_control.service
[Unit]
Description=R2 Motor Control Web Server
After=pigpiod.service pre-cache.service
Requires=pigpiod.service
Wants=pre-cache.service

[Service]
Type=simple
User=r2tele
WorkingDirectory=/home/r2tele/motor_control
ExecStart=/usr/bin/python3 /home/r2tele/motor_control/app.py
Restart=on-failure
RestartSec=2
StartLimitInterval=120
StartLimitBurst=5
MemoryMax=384M
LimitNOFILE=1024

[Install]
WantedBy=multi-user.target
